# linux/aosp/aosp.service
[Unit]
Description=aosp master daily build service

[Service]
Type=simple
User=dp
ExecStart=/home/dp/code/droid/linux/aosp/makeaosp.sh
